# deploy/codeburg.service
[Unit]
Description=Codeburg Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=codeburg
Group=codeburg
WorkingDirectory=/opt/codeburg
ExecStart=/opt/codeburg/backend/codeburg serve --host 127.0.0.1 --port 8080
Restart=on-failure
RestartSec=5
KillMode=process
Environment=HOME=/home/codeburg
Environment=TERM=xterm-256color

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/codeburg
ReadWritePaths=/opt/codeburg
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
